# images/desktop/config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── System D-Bus daemon (budgie-session refuses to start without one) ───────
[program:dbus]
command=dbus-daemon --system --nofork --nopidfile
autorestart=true
priority=100
startsecs=1
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus.log

# ── Budgie desktop session (labwc, headless Wayland — see Dockerfile) ───────
# WLR_BACKENDS/WLR_LIBINPUT_NO_DEVICES/WAYLAND_DISPLAY/XDG_RUNTIME_DIR come
# from the image's ENV. startbudgielabwc is Budgie's own launcher (ships as
# the Exec= of /usr/share/wayland-sessions/budgie-desktop.desktop) — it wires
# up labwc as the compositor and starts budgie-desktop as its session client.
[program:budgie]
command=startbudgielabwc
user=%(ENV_USER)s
# XDG_CURRENT_DESKTOP=Budgie:GNOME (not just "Budgie"): budgie-session-binary
# sets it to bare "Budgie" itself, which makes Budgie's own menu hide every
# .desktop file marked OnlyShowIn=GNOME (gnome-terminal, and most of GNOME's
# own apps) even though they run fine — this is what "some apps don't show
# up/won't launch from the menu" turned out to be. Matches upstream Ubuntu
# Budgie's convention.
environment=HOME=%(ENV_HOME)s,USER=%(ENV_USER)s,XDG_RUNTIME_DIR=%(ENV_XDG_RUNTIME_DIR)s,WLR_BACKENDS=%(ENV_WLR_BACKENDS)s,WLR_LIBINPUT_NO_DEVICES=%(ENV_WLR_LIBINPUT_NO_DEVICES)s,XDG_CURRENT_DESKTOP="Budgie:GNOME"
autorestart=true
priority=200
startsecs=3
stdout_logfile=/var/log/supervisor/budgie.log
stderr_logfile=/var/log/supervisor/budgie.log

# ── wayvnc: captures the labwc output over wlr-screencopy ────────────────────
# Waits for the compositor's Wayland socket instead of a flat sleep, same
# reasoning as the old X11 readiness poll this replaces.
[program:wayvnc]
command=bash -c "until [ -S $XDG_RUNTIME_DIR/%(ENV_WAYLAND_DISPLAY)s ]; do sleep 0.1; done; exec wayvnc 127.0.0.1 %(ENV_VNC_PORT)s"
user=%(ENV_USER)s
environment=HOME=%(ENV_HOME)s,USER=%(ENV_USER)s,XDG_RUNTIME_DIR=%(ENV_XDG_RUNTIME_DIR)s,WAYLAND_DISPLAY=%(ENV_WAYLAND_DISPLAY)s
autorestart=true
priority=300
startsecs=5
stdout_logfile=/var/log/supervisor/wayvnc.log
stderr_logfile=/var/log/supervisor/wayvnc.log

# ── websockify: noVNC static files + WebSocket→VNC bridge ────────────────────
[program:websockify]
command=/usr/local/bin/websockify-fork
    --web /opt/novnc
    --heartbeat=30
    0.0.0.0:%(ENV_PROXY_PORT)s
    127.0.0.1:%(ENV_VNC_PORT)s
autorestart=true
priority=400
startsecs=5
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log
