# images/desktop/config/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── Virtual display ───────────────────────────────────────────────────────────
[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_VNC_RESOLUTION)sx%(ENV_VNC_DEPTH)s -ac -nolisten tcp
user=%(ENV_USER)s
autorestart=true
priority=100
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.log

# ── Disable screensaver/DPMS (runs once after Xvfb) ──────────────────────────
[program:disable-screensaver]
command=bash -c "sleep 3 && xset s off -dpms s noblank"
user=%(ENV_USER)s
environment=DISPLAY=%(ENV_DISPLAY)s,HOME=%(ENV_HOME)s
autostart=true
autorestart=false
priority=150
stdout_logfile=/dev/null
stderr_logfile=/dev/null

# ── XFCE4 desktop session ─────────────────────────────────────────────────────
[program:xfce4]
command=startxfce4
user=%(ENV_USER)s
environment=DISPLAY=%(ENV_DISPLAY)s,HOME=%(ENV_HOME)s,USER=%(ENV_USER)s
autorestart=true
priority=200
startsecs=3
stdout_logfile=/var/log/supervisor/xfce4.log
stderr_logfile=/var/log/supervisor/xfce4.log

# ── x11vnc VNC server (localhost only) ───────────────────────────────────────
# -wait 5 / -defer 5: poll every 5ms, batch updates within 5ms → lower latency
# -threads: parallel encoding
# XDamage and XFIXES are left enabled (default) for correct cursor sprites
[program:x11vnc]
command=x11vnc
    -display %(ENV_DISPLAY)s
    -rfbport %(ENV_VNC_PORT)s
    -nopw
    -forever
    -shared
    -localhost
    -wait 5
    -defer 5
    -threads
    -repeat
user=%(ENV_USER)s
environment=DISPLAY=%(ENV_DISPLAY)s,HOME=%(ENV_HOME)s,USER=%(ENV_USER)s
autorestart=true
priority=300
startsecs=5
stdout_logfile=/var/log/supervisor/x11vnc.log
stderr_logfile=/var/log/supervisor/x11vnc.log

# ── websockify: noVNC static files + WebSocket→VNC bridge ────────────────────
[program:websockify]
command=/usr/local/bin/websockify-fork
    --web /opt/novnc
    --heartbeat=30
    0.0.0.0:%(ENV_PROXY_PORT)s
    127.0.0.1:%(ENV_VNC_PORT)s
autorestart=true
priority=400
startsecs=5
stdout_logfile=/var/log/supervisor/websockify.log
stderr_logfile=/var/log/supervisor/websockify.log
